# src/influxdb/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:influxdb]
command=influxd -config /etc/influxdb.conf

[program:telegraf]
command=/usr/bin/telegraf -config /etc/telegraf/telegraf.conf
