# src/php_my_admin/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
[program:php7.3-fpm]
command=/usr/sbin/php-fpm7 --nodaemonize
user=root
autostart=true
autorestart=true