# src/grafana/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:grafana]
#command=./grafana-server -config "grafana-7.5.4/conf/defaults.ini" -homepath "/grafana-7.5.4"
command=./grafana-server
#command=sh
