# src/mysql/supervisor.conf
[supervisord]
nodaemon=true
user=root

[program:telegraf]
command = sh
autorestart = true

[program:mysql]
command = /usr/bin/mysqld_safe --datadir="/var/lib/mysql"
autorestart = true